# scripts/com.antigravity.ot_automation.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.antigravity.ot_automation</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kk/Desktop/오티 자동화/scripts/keep_alive.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kk/Desktop/오티 자동화/outputs/logs/launch_agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kk/Desktop/오티 자동화/outputs/logs/launch_agent.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/kk/Desktop/오티 자동화</string>
</dict>
</plist>
